# flaskapp.service
[Unit]
Description=Flask App - Sistema de Formularios
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/flaskapp
Environment="PATH=/var/www/flaskapp/venv/bin"
Environment="FLASK_ENV=production"
# Puedes sobreescribir SECRET_KEY creando /etc/flaskapp.env
EnvironmentFile=-/etc/flaskapp.env
ExecStart=/var/www/flaskapp/venv/bin/gunicorn --config gunicorn.conf.py app:app
ExecStop=/bin/kill -s TERM $MAINPID
KillMode=mixed
TimeoutStopSec=5
PrivateTmp=true
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
